# test_pkg/etc/systemd/system/time_script2.service
[Unit]
Description=Time Script Service

[Service]
ExecStart=/usr/sbin/time_script.sh

[Install]
WantedBy=multi-user.target
